# deploy/systemd/botonio-botsci@staging.service.d/override.conf
# Staging instance of botonio-botsci@.service - wired to the TEST-server Discord application
# with non-production tokens. Fill the <placeholders> on the box (Task 13). The credential
# blobs are produced with `systemd-creds encrypt`.
[Service]
User=botonio-botsci-staging
Group=botonio-botsci-staging

# Secrets (delivered into $CREDENTIALS_DIRECTORY, never the process environment).
# One Discord bot token serves both the gateway and the HTTP API (read as DISCORD_BOT_TOKEN).
LoadCredentialEncrypted=discord_bot_token:/etc/botonio-botsci/staging/discord_bot_token.cred
LoadCredentialEncrypted=solidarity_tech_token:/etc/botonio-botsci/staging/solidarity_tech_token.cred

# Non-secret config (the test server + application).
Environment=DISCORD_GUILD_ID=<staging-guild-id>
Environment=DISCORD_MODERATOR_ROLE_ID=<staging-moderator-role-id>
Environment=SOLIDARITY_TECH_DISCORD_LIST_ID=<staging-user-list-id>
Environment=BOT_ACCENT_COLOR=#c8102e
Environment=BOT_INDEX_REFRESH_SECS=900
Environment=RUST_LOG=info,botonio_botsci=debug

# Divorced staging: read fabricated members from an in-process mock Solidarity Tech server
# instead of the real API. SOLIDARITY_TECH_MOCK enables it; the bot binds the mock to the
# host:port of SOLIDARITY_TECH_BASE_URL (the same URL the client reads), so the two cannot
# drift. SOLIDARITY_TECH_MOCK_PERSONAS binds real test-server Discord account ids to
# fabricated member states - fill the <id> placeholders on the box. Personas: good_standing,
# amber, lapsed, retired_tier, malformed (an unmapped account reads as "not a member").
# Leave all three unset to read the real API.
Environment=SOLIDARITY_TECH_MOCK=1
Environment=SOLIDARITY_TECH_BASE_URL=http://127.0.0.1:8080
Environment=SOLIDARITY_TECH_MOCK_PERSONAS=<id>=good_standing,<id>=amber,<id>=lapsed
