# deploy/systemd/botonio-botsci@staging.service.d/override.conf.tmpl
# Staging instance of botonio-botsci@.service - wired to the TEST-server Discord application
# with non-production tokens. Fill the <placeholders> on the box before starting. The credential
# blobs are produced with `systemd-creds encrypt`.
[Service]
User=botonio-botsci-staging
Group=botonio-botsci-staging

# Secrets (delivered into $$CREDENTIALS_DIRECTORY, never the process environment).
# One Discord bot token serves both the gateway and the HTTP API (read as DISCORD_BOT_TOKEN).
LoadCredentialEncrypted=discord_bot_token:/etc/botonio-botsci/bot/staging/discord_bot_token.cred
# The HMAC key that hashes the Discord ids written to the append-only audit log.
LoadCredentialEncrypted=audit_hash_key:/etc/botonio-botsci/bot/staging/audit_hash_key.cred

# Non-secret config (the test server + application).
Environment=DISCORD_GUILD_ID=${discord_guild_id}
Environment=SOLIDARITY_TECH_DISCORD_LIST_ID=${st_user_list_id}
Environment=BOT_ACCENT_COLOR=#c8102e
Environment=BOT_INDEX_REFRESH_SECS=900
Environment=RUST_LOG=info,botonio_botsci=debug
# Audit log (the append-only record of moderator card lookups): AUDIT_KEY_ID is stamped on
# each row and bumped when audit_hash_key is rotated; BOT_LOOKUP_RATE_PER_MIN caps how many
# other-member card lookups one moderator may run per minute.
Environment=AUDIT_KEY_ID=v1
Environment=BOT_LOOKUP_RATE_PER_MIN=10

# Scheduled membership scan (role sync + orphan sweep). OFF until enabled per-guild in
# /setup (scan_enabled); these only tune it once it is on. Staging uses the short floor
# interval for a quick test loop; the tripwire aborts a pass that would demote at least
# PERCENT% of scanned members AND at least FLOOR members.
Environment=BOT_SCAN_INTERVAL_SECS=300
Environment=BOT_SCAN_TRIPWIRE_PERCENT=20
Environment=BOT_SCAN_TRIPWIRE_FLOOR=5
Environment=BOT_SCAN_PACE_MS=250

# SSO role-check endpoint. The three credentials are delivered into $$CREDENTIALS_DIRECTORY
# (never the process environment) and loaded by the bot on startup when BOT_SSO_ENABLED=1.
# Leave BOT_SSO_ENABLED unset until workspace-sync is deployed and the three real credentials
# are provisioned; with it unset the socket never binds and the SSO path is entirely inert.
#
# To enable: create the botonio-staging-sso group with the bot's and workspace-sync's service
# users as its only members, set BOT_SSO_ENABLED=1 here, set sso_oauth_client_id and
# sso_redirect_uri in secrets/staging.enc.yaml (they render into the lines below, like
# discord_guild_id), and provision the three .cred files. botonio-staging-sso is a
# SUPPLEMENTARY group: the bot's primary group stays its own, so handing the socket to this
# group never exposes the bot's credential directory to workspace-sync.
SupplementaryGroups=botonio-staging-sso
LoadCredentialEncrypted=sso_oauth_client_secret:/etc/botonio-botsci/bot/staging/sso_oauth_client_secret.cred
LoadCredentialEncrypted=sso_signing_key:/etc/botonio-botsci/bot/staging/sso_signing_key.cred
LoadCredentialEncrypted=sso_caller_bearer:/etc/botonio-botsci/bot/staging/sso_caller_bearer.cred
Environment=BOT_SSO_ENABLED=1
Environment=BOT_SSO_SOCKET_PATH=/run/botonio-staging/sso.sock
Environment=BOT_SSO_SOCKET_GROUP=botonio-staging-sso
Environment=BOT_SSO_AUDIENCE=workspace-sync
Environment=BOT_SSO_REDIRECT_URI=${sso_redirect_uri}
Environment=BOT_SSO_OAUTH_CLIENT_ID=${sso_oauth_client_id}
Environment=BOT_SSO_KEY_ID=v1

# Local membership cache (Postgres). The runtime connection authenticates by peer over
# the Unix socket - no password (the long-running process holds no DB secret). The
# migrate phase (ExecStartPre) uses the migration credential below over TCP loopback.
LoadCredentialEncrypted=db_migration_password:/etc/botonio-botsci/bot/staging/db_migration_password.cred
Environment=DB_NAME=botonio_staging
Environment=DB_RUNTIME_DSN=postgres:///botonio_staging?host=/var/run/postgresql&user=botonio_staging_app
# The migrate phase connects to the botonio_staging_migrate role over TCP loopback, default
# 127.0.0.1:5432. Uncomment to point it at a cluster that listens elsewhere:
# Environment=DB_MIGRATE_HOST=127.0.0.1
# Environment=DB_MIGRATE_PORT=5432

# Divorced staging: read fabricated members from an in-process mock, not the real API.
# SOLIDARITY_TECH_MOCK enables it; the mock binds the host:port of SOLIDARITY_TECH_BASE_URL
# (the URL the client reads, so they can't drift); _PERSONAS maps account ids to fabricated
# states (good_standing/amber/lapsed/retired_tier/malformed; an unmapped id reads as unknown).
# Unset all three for the real API. Pick a free loopback port for the bind - 8080 is a poor
# choice on a box already running crowdsec or a similar service there.
Environment=SOLIDARITY_TECH_MOCK=1
Environment=SOLIDARITY_TECH_BASE_URL=http://127.0.0.1:9721
Environment=SOLIDARITY_TECH_MOCK_PERSONAS=${good_standing_id}=good_standing,${expiring_standing_id}=amber,${lapsed_user_id}=lapsed,test@example.com=email_verify
# Dummy token (the mock ignores it): staging holds no real Solidarity Tech credential, so a
# staging compromise can't reach real member PII - and with only this, real ST 401s anyway.
Environment=SOLIDARITY_TECH_TOKEN=mock

# Staging-only testing affordances, each gated on its own env var (absent in production):
# BOT_FORGET_COMMAND registers /forget (reset one member); BOT_STRIP_ROLES_COMMAND registers
# /strip-roles (clear the managed roles server-wide so /bulk-verify can be retested). The
# forget-of-overridden members /strip-roles runs also relies on the staging DELETE grant.
Environment=BOT_FORGET_COMMAND=1
Environment=BOT_STRIP_ROLES_COMMAND=1

# BOT_REMINDER_TEST_SEND registers /reminders test-send (fire one dues-reminder milestone at a
# member on demand) so the reminder render/thread/button path can be walked with one test
# account. Staging only - production never registers it.
Environment=BOT_REMINDER_TEST_SEND=1
