# deploy/systemd/botonio-botsci@staging.service.d/override.conf
# Staging instance of botonio-botsci@.service - wired to the TEST-server Discord application
# with non-production tokens. Fill the <placeholders> on the box (Task 13). The credential
# blobs are produced with `systemd-creds encrypt`.
[Service]
User=botonio-botsci-staging
Group=botonio-botsci-staging

# Secrets (delivered into $CREDENTIALS_DIRECTORY, never the process environment).
# One Discord bot token serves both the gateway and the HTTP API (read as DISCORD_BOT_TOKEN).
LoadCredentialEncrypted=discord_bot_token:/etc/botonio-botsci/staging/discord_bot_token.cred
LoadCredentialEncrypted=solidarity_tech_token:/etc/botonio-botsci/staging/solidarity_tech_token.cred

# Non-secret config (the test server + application).
Environment=DISCORD_GUILD_ID=<staging-guild-id>
Environment=DISCORD_MODERATOR_ROLE_ID=<staging-moderator-role-id>
Environment=SOLIDARITY_TECH_DISCORD_LIST_ID=<staging-user-list-id>
Environment=BOT_ACCENT_COLOR=#c8102e
Environment=BOT_INDEX_REFRESH_SECS=900
Environment=RUST_LOG=info,botonio_botsci=debug

# Divorced-staging seam (enable with the mock-Solidarity-Tech slice): point the Solidarity
# Tech client at a local mock so staging reads no real member records. Until then staging
# uses the real API.
# Environment=SOLIDARITY_TECH_BASE_URL=http://127.0.0.1:8080
