# deploy/systemd/bot-db-backup.timer
# Nightly trigger for bot-db-backup.service. Enable with:
#   systemctl enable --now bot-db-backup.timer
[Unit]
Description=Nightly encrypted off-box backup of the production membership database

[Timer]
OnCalendar=*-*-* 03:30:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
