# scripts/botonio_deploy/assets/botonio-botsci@production.service.d/override.conf.tmpl
# Production instance of botonio-botsci@.service - the REAL Discord application + server.
# The credential blobs are produced with `systemd-creds encrypt`.
[Service]
User=botonio-botsci-production
Group=botonio-botsci-production

# Secrets (delivered into $$CREDENTIALS_DIRECTORY, never the process environment).
# One Discord bot token serves both the gateway and the HTTP API (read as DISCORD_BOT_TOKEN).
LoadCredentialEncrypted=discord_bot_token:/etc/botonio-botsci/bot/production/discord_bot_token.cred
LoadCredentialEncrypted=solidarity_tech_token:/etc/botonio-botsci/bot/production/solidarity_tech_token.cred
# The HMAC key that hashes the Discord ids written to the append-only audit log.
LoadCredentialEncrypted=audit_hash_key:/etc/botonio-botsci/bot/production/audit_hash_key.cred

# Non-secret config (the real server + application).
Environment=DISCORD_GUILD_ID=${discord_guild_id}
Environment=SOLIDARITY_TECH_DISCORD_LIST_ID=${st_user_list_id}
Environment=BOT_ACCENT_COLOR=#c8102e
Environment=BOT_INDEX_REFRESH_SECS=900
Environment=RUST_LOG=info,botonio_botsci=info
# Audit log (the append-only record of moderator card lookups): AUDIT_KEY_ID is stamped on
# each row and bumped when audit_hash_key is rotated; BOT_LOOKUP_RATE_PER_MIN caps how many
# other-member card lookups one moderator may run per minute.
Environment=AUDIT_KEY_ID=v1
Environment=BOT_LOOKUP_RATE_PER_MIN=10

# Local membership cache (Postgres). The runtime connection authenticates by peer over
# the Unix socket - no password. The migrate phase (ExecStartPre) uses the migration
# credential below over TCP loopback to apply the schema.
LoadCredentialEncrypted=db_migration_password:/etc/botonio-botsci/bot/production/db_migration_password.cred
Environment=DB_NAME=botonio_production
Environment=DB_RUNTIME_DSN=postgres:///botonio_production?host=/var/run/postgresql&user=botonio_production_app
# The migrate phase connects to the botonio_production_migrate role over TCP loopback, default
# 127.0.0.1:5432. Uncomment to point it at a cluster that listens elsewhere:
# Environment=DB_MIGRATE_HOST=127.0.0.1
# Environment=DB_MIGRATE_PORT=5432
