# deploy/systemd/botonio-botsci@production.service.d/override.conf.tmpl
# Production instance of botonio-botsci@.service - the REAL Discord application + server.
# The credential blobs are produced with `systemd-creds encrypt`.
[Service]
User=botonio-botsci-production
Group=botonio-botsci-production

# Secrets (delivered into $$CREDENTIALS_DIRECTORY, never the process environment).
# One Discord bot token serves both the gateway and the HTTP API (read as DISCORD_BOT_TOKEN).
LoadCredentialEncrypted=discord_bot_token:/etc/botonio-botsci/bot/production/discord_bot_token.cred
LoadCredentialEncrypted=solidarity_tech_token:/etc/botonio-botsci/bot/production/solidarity_tech_token.cred
# The HMAC key that hashes the Discord ids written to the append-only audit log.
LoadCredentialEncrypted=audit_hash_key:/etc/botonio-botsci/bot/production/audit_hash_key.cred

# Non-secret config (the real server + application).
Environment=DISCORD_GUILD_ID=${discord_guild_id}
Environment=SOLIDARITY_TECH_DISCORD_LIST_ID=${st_user_list_id}
Environment=BOT_ACCENT_COLOR=#c8102e
Environment=BOT_INDEX_REFRESH_SECS=900
Environment=RUST_LOG=info,botonio_botsci=info
# Audit log (the append-only record of moderator card lookups): AUDIT_KEY_ID is stamped on
# each row and bumped when audit_hash_key is rotated; BOT_LOOKUP_RATE_PER_MIN caps how many
# other-member card lookups one moderator may run per minute.
Environment=AUDIT_KEY_ID=v1
Environment=BOT_LOOKUP_RATE_PER_MIN=10

# Scheduled membership scan (role sync + orphan sweep). OFF until enabled per-guild in
# /setup (scan_enabled); these only tune it once it is on. 14400s = 4h matches Solidarity
# Tech's own record-sync cadence; the tripwire aborts a pass that would demote at least
# PERCENT% of scanned members AND at least FLOOR members.
Environment=BOT_SCAN_INTERVAL_SECS=14400
Environment=BOT_SCAN_TRIPWIRE_PERCENT=20
Environment=BOT_SCAN_TRIPWIRE_FLOOR=5
Environment=BOT_SCAN_PACE_MS=250

# SSO role-check endpoint. The three credentials are delivered into $$CREDENTIALS_DIRECTORY
# (never the process environment) and loaded by the bot on startup when BOT_SSO_ENABLED=1.
# Leave BOT_SSO_ENABLED unset until workspace-sync is deployed and the three real credentials
# are provisioned; with it unset the socket never binds and the SSO path is entirely inert.
#
# To enable: create the botonio-production-sso group with the bot's and workspace-sync's
# service users as its only members, set BOT_SSO_ENABLED=1 here, set sso_oauth_client_id and
# sso_redirect_uri in secrets/production.enc.yaml (they render into the lines below, like
# discord_guild_id), and provision the three .cred files. botonio-production-sso is a
# SUPPLEMENTARY group: the bot's primary group stays its own, so handing the socket to this
# group never exposes the bot's credential directory to workspace-sync.
SupplementaryGroups=botonio-production-sso
LoadCredentialEncrypted=sso_oauth_client_secret:/etc/botonio-botsci/bot/production/sso_oauth_client_secret.cred
LoadCredentialEncrypted=sso_signing_key:/etc/botonio-botsci/bot/production/sso_signing_key.cred
LoadCredentialEncrypted=sso_caller_bearer:/etc/botonio-botsci/bot/production/sso_caller_bearer.cred
# Environment=BOT_SSO_ENABLED=1 # TODO: ENABLE THIS DUMMY WHEN WORKSPACE SYNC IS READY TO HIT PROD
Environment=BOT_SSO_SOCKET_PATH=/run/botonio-production/sso.sock
Environment=BOT_SSO_SOCKET_GROUP=botonio-production-sso
Environment=BOT_SSO_AUDIENCE=workspace-sync
Environment=BOT_SSO_REDIRECT_URI=${sso_redirect_uri}
Environment=BOT_SSO_OAUTH_CLIENT_ID=${sso_oauth_client_id}
Environment=BOT_SSO_KEY_ID=v1

# Local membership cache (Postgres). The runtime connection authenticates by peer over
# the Unix socket - no password. The migrate phase (ExecStartPre) uses the migration
# credential below over TCP loopback to apply the schema.
LoadCredentialEncrypted=db_migration_password:/etc/botonio-botsci/bot/production/db_migration_password.cred
Environment=DB_NAME=botonio_production
Environment=DB_RUNTIME_DSN=postgres:///botonio_production?host=/var/run/postgresql&user=botonio_production_app
# The migrate phase connects to the botonio_production_migrate role over TCP loopback, default
# 127.0.0.1:5432. Uncomment to point it at a cluster that listens elsewhere:
# Environment=DB_MIGRATE_HOST=127.0.0.1
# Environment=DB_MIGRATE_PORT=5432
