# deploy/systemd/botonio-botsci@production.service.d/override.conf
# Production instance of botonio-botsci@.service - the REAL Discord application + server. 
# The credential blobs are produced with `systemd-creds encrypt`.
[Service]
User=botonio-botsci-production
Group=botonio-botsci-production

# Secrets (delivered into $CREDENTIALS_DIRECTORY, never the process environment).
# One Discord bot token serves both the gateway and the HTTP API (read as DISCORD_BOT_TOKEN).
LoadCredentialEncrypted=discord_bot_token:/etc/botonio-botsci/production/discord_bot_token.cred
LoadCredentialEncrypted=solidarity_tech_token:/etc/botonio-botsci/production/solidarity_tech_token.cred

# Non-secret config (the real server + application).
Environment=DISCORD_GUILD_ID=<production-guild-id>
Environment=DISCORD_MODERATOR_ROLE_ID=<production-moderator-role-id>
Environment=SOLIDARITY_TECH_DISCORD_LIST_ID=<production-user-list-id>
Environment=BOT_ACCENT_COLOR=#c8102e
Environment=BOT_INDEX_REFRESH_SECS=900
Environment=RUST_LOG=info,botonio_botsci=info
