# scripts/botonio_deploy/assets/botonio-botsci@staging.service.d/override.conf.tmpl
# Staging instance of botonio-botsci@.service - wired to the TEST-server Discord application
# with non-production tokens. Fill the <placeholders> on the box before starting. The credential
# blobs are produced with `systemd-creds encrypt`.
[Service]
User=botonio-botsci-staging
Group=botonio-botsci-staging

# Secrets (delivered into $$CREDENTIALS_DIRECTORY, never the process environment).
# One Discord bot token serves both the gateway and the HTTP API (read as DISCORD_BOT_TOKEN).
LoadCredentialEncrypted=discord_bot_token:/etc/botonio-botsci/bot/staging/discord_bot_token.cred
# The HMAC key that hashes the Discord ids written to the append-only audit log.
LoadCredentialEncrypted=audit_hash_key:/etc/botonio-botsci/bot/staging/audit_hash_key.cred

# Non-secret config (the test server + application).
Environment=DISCORD_GUILD_ID=${discord_guild_id}
Environment=SOLIDARITY_TECH_DISCORD_LIST_ID=${st_user_list_id}
Environment=BOT_ACCENT_COLOR=#c8102e
Environment=BOT_INDEX_REFRESH_SECS=900
Environment=RUST_LOG=info,botonio_botsci=debug
# Audit log (the append-only record of moderator card lookups): AUDIT_KEY_ID is stamped on
# each row and bumped when audit_hash_key is rotated; BOT_LOOKUP_RATE_PER_MIN caps how many
# other-member card lookups one moderator may run per minute.
Environment=AUDIT_KEY_ID=v1
Environment=BOT_LOOKUP_RATE_PER_MIN=10

# Local membership cache (Postgres). The runtime connection authenticates by peer over
# the Unix socket - no password (the long-running process holds no DB secret). The
# migrate phase (ExecStartPre) uses the migration credential below over TCP loopback.
LoadCredentialEncrypted=db_migration_password:/etc/botonio-botsci/bot/staging/db_migration_password.cred
Environment=DB_NAME=botonio_staging
Environment=DB_RUNTIME_DSN=postgres:///botonio_staging?host=/var/run/postgresql&user=botonio_staging_app
# The migrate phase connects to the botonio_staging_migrate role over TCP loopback, default
# 127.0.0.1:5432. Uncomment to point it at a cluster that listens elsewhere:
# Environment=DB_MIGRATE_HOST=127.0.0.1
# Environment=DB_MIGRATE_PORT=5432

# Divorced staging: read fabricated members from an in-process mock, not the real API.
# SOLIDARITY_TECH_MOCK enables it; the mock binds the host:port of SOLIDARITY_TECH_BASE_URL
# (the URL the client reads, so they can't drift); _PERSONAS maps account ids to fabricated
# states (good_standing/amber/lapsed/retired_tier/malformed; an unmapped id reads as unknown).
# Unset all three for the real API. Pick a free loopback port for the bind - 8080 is a poor
# choice on a box already running crowdsec or a similar service there.
Environment=SOLIDARITY_TECH_MOCK=1
Environment=SOLIDARITY_TECH_BASE_URL=http://127.0.0.1:9721
Environment=SOLIDARITY_TECH_MOCK_PERSONAS=${good_standing_id}=good_standing,${expiring_standing_id}=amber,${lapsed_user_id}=lapsed,test@example.com=email_verify
# Dummy token (the mock ignores it): staging holds no real Solidarity Tech credential, so a
# staging compromise can't reach real member PII - and with only this, real ST 401s anyway.
Environment=SOLIDARITY_TECH_TOKEN=mock

# Staging-only testing affordances, each gated on its own env var (absent in production):
# BOT_FORGET_COMMAND registers /forget (reset one member); BOT_STRIP_ROLES_COMMAND registers
# /strip-roles (clear the managed roles server-wide so /bulk-verify can be retested). The
# forget-of-overridden members /strip-roles runs also relies on the staging DELETE grant.
Environment=BOT_FORGET_COMMAND=1
Environment=BOT_STRIP_ROLES_COMMAND=1
